# systemd/user/restic-backup.timer
[Unit]
Description=Timer for restic-backup

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
